# provisioning/etc/systemd/system/sevi-control.socket
[Unit]
Description=SeviControl socket

[Socket]
ListenStream=/run/sevi-control/socket

[Install]
WantedBy=sockets.target
